# confs/rendered/carebot.bot.conf
# description "carebot app"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
    echo "Welcome" >> /home/ubuntu/apps/carebot/log.txt

    bash -c '/usr/local/bin/virtualenvwrapper.sh'
    echo "Wrapper created" >> /home/ubuntu/apps/carebot/log.txt

    # . /home/ubuntu/apps/carebot/.env
    # echo "Env variables sourced" >> /home/ubuntu/apps/carebot/log.txt

    . /home/ubuntu/apps/carebot/virtualenv/bin/activate
    echo "Env activated" >> /home/ubuntu/apps/carebot/log.txt

    cd /home/ubuntu/apps/carebot
    echo "Changed to project path" >> /home/ubuntu/apps/carebot/log.txt

    exec virtualenv/bin/python carebot.py >> /home/ubuntu/apps/carebot/log.txt
end script